# supervisord.conf
[supervisord]
nodaemon=true
# https://stackoverflow.com/questions/63608075/userwarning-supervisord-is-running-as-root-and-it-is-searching-for-its-configur
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:freeswitch]
autorestart=true
command=/bin/sh /usr/bin/start-freeswitch.sh
